# bjtu-needs-auth.target
[Unit]
Description=BJTU Network Authentication Required